# myControlScript.sh
#!/bin/bash
# chkconfig: 2345 1 80
# This will begin my python code, I hope :)

#/usr/CRT_Comms/myScript.sh &
/usr/CRT_Comms/mySScript.sh &
